# crontab
*/10 * * * * cd /app && docker compose run --rm send-photo